# bin/cron/tasks/default.cron
#crontab
* * * * *   php ../../php/backup.php # aaaa
